# ansible/roles/initial/files/nginx/cron.d/nginx_count_access
1-59/5 * * * * root /usr/local/bin/nginx_count_access.py | logger -i -p cron.info -t 'nginx_count_access'
